# systemd/cursor-update.timer
[Unit]
Description=Daily Cursor IDE update

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

