# systemd/cursor-update.service
[Unit]
Description=Update Cursor IDE (AppImage)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-cursor

